# cosmos.service
[Unit]
Description=Cosmos
After=multi-user.target

[Service]
Type=simple
WorkingDirectory=/root/cosmos-discord-bot
ExecStart=/usr/bin/python3.6 /root/cosmos-discord-bot/run.py
User=root
Restart=always

[Install]
WantedBy=multi-user.target
